# scripts/assist-web.service.template
[Unit]
Description=Assist Web Service
After=network.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{DEPLOY_PATH}}
Environment="PATH={{DEPLOY_PATH}}/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="ASSIST_THREADS_DIR={{ASSIST_THREADS_DIR}}"
{{ENVIRONMENT_VARS}}
ExecStart={{DEPLOY_PATH}}/.venv/bin/python -m manage.web
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assist-web

[Install]
WantedBy=multi-user.target
